# launchd/com.prcontroller.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template — optional standalone poller (useful when the server is not running).
  Run scripts/install.sh --with-poller to activate.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prcontroller.poll</string>

  <!-- Run every 3 minutes -->
  <key>StartInterval</key>
  <integer>180</integer>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>pr_controller</string>
    <string>poll</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PRCTL_CONFIG</key>
    <string>__HOME__/.pr-controller/config.yaml</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/.pr-controller/poll.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.pr-controller/poll.log</string>
</dict>
</plist>
